# crontab
*/15 * * * * root /var/scripts/sync.sh >> /var/log/cron.log 2>&1
